# .config/autostart/redshift-gtk.desktop
[Desktop Entry]
Name[es]=Redshift
GenericName[hu]=Szín hőmérséklet beállítás
Exec=redshift-gtk
X-GNOME-Autostart-enabled=false
Comment[tr]=Renk sıcaklığı ayarlama aracı
Comment[cs]=Nástroj pro přizpůsobení barevné teploty
Name[zh_CN]=红移
Comment[it]=Strumento per la regolazione della temperatura del colore
Icon=redshift
GenericName[fr]=Réglage de la température de couleur
GenericName[lt]=Spalvos temperatūros reguliavimas
StartupNotify=true
GenericName[tr]=Renk sıcaklığı ayarı
Comment[lt]=Spalvos temperatūros reguliavimo įrankis
GenericName[pl]=Dostosowanie temperatury barwowej
GenericName[cs]=Přizpůsobení barevné teploty
Type=Application
Version=1.0
Comment[pl]=Narzędzie do dostosowywania temperatury barwowej
Comment[hu]=Színhőmérséklet beállító eszköz
Comment[zh_CN]=色温调节工具
Comment[nl]=Kleurtemperatuur aanpassingsprogramma
Name[it]=Redshift
Terminal=false
Name[fr]=Redshift
Comment=Color temperature adjustment tool
Hidden=true
Comment[fr]=Outil de réglage de la température de couleur
Comment[ru]=Инструмент регулирования цветовой температуры
Comment[da]=Justeringsværktøj for farvetemperatur
Name[ru]=Redshift
GenericName[zh_CN]=色温调节
GenericName=Color temperature adjustment
Categories=Utility;
Comment[es]=Herramienta para el ajuste de la temperatura de color
Name[lt]=Redshift
GenericName[de]=Farbtemperatur Anpassung
GenericName[es]=Ajuste de la temperatura de color
Name[tr]=Redshift
Name[pl]=Redshift
GenericName[nl]=Kleurtemperatuur aanpassing
Name[de]=Redshift
Name=Redshift
GenericName[ru]=Настройка цветовой температуры
GenericName[da]=Justering af farvetemperatur
Comment[de]=Farbtemperatur Anpassungswerkzeug
Name[pt_BR]=Redshift
Name[hu]=Redshift
GenericName[it]=Regolazione della temperatura del colore

